# rc.local
#!/bin/sh -e
#
# rc.local
#
# prototype rc.local script for pi-cv-env
#

chmod 777 /tmp
chmod +t  /tmp

if [ -x /home/pi/example/start.sh ]
then
  su pi -c "/home/pi/examples/start.sh >/tmp/example.log 2>&1"
fi

exit 0

